# network/if-up.d/static-arp
#! / Bin / sh
arp -i eth0 -s 193.191.177.254 0A: 30: 5A: 4D: 0B: D0
